# todo-hotel.service
[Unit]
Description=Todo Hotel Application
Documentation=file:///root/todo-hotel/DEPLOIEMENT_SIMPLE.md
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/todo-hotel
Environment=PATH=/root/todo-hotel/venv/bin
Environment=FLASK_APP=app:create_app
Environment=FLASK_ENV=production
# Bind on 0.0.0.0:8081 for Tailscale access (no Nginx)
ExecStart=/root/todo-hotel/venv/bin/gunicorn -w 3 -b 0.0.0.0:8081 wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=10
# Simple, no advanced sandboxing for root
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Load environment variables from .env if present
EnvironmentFile=-/root/todo-hotel/.env

[Install]
WantedBy=multi-user.target
